# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:kingdom]
command=/usr/bin/python kingdom.py
directory=/deploy/main/app
